# sv/pipewire/run
#!/bin/sh

exec 2>&1
sv start dbus-session # Get dbus-session if we can, like dinit(8) <waits-for>
set -efu
. "$SVDIR"/../conf.sh

exec pipewire
