# sv/iwgtk/run
#!/bin/sh
exec 2>&1
set -efu
. "$SVDIR"/../conf.sh
sv check wayland
. "$XDG_RUNTIME_DIR"/swayenv.$(cat "$SVDIR"/sway/supervise/pid)
exec iwgtk -i
