# sv/dbus-session/run
#!/bin/sh
set -efu
. "$SVDIR"/../conf.sh # So dbus-started processes inherit SSH_AUTH_SOCK, and we get dbus_session_bus
exec dbus-daemon --session --nofork --nopidfile --syslog-only --address=unix:path="$dbus_session_bus"
