# sv/ssh-agent/run
#!/bin/sh
exec ssh-agent -a "$XDG_RUNTIME_DIR/ssh-agent" -t 3h -D 2>&1
