# sv/wireplumber/run
#!/bin/sh

exec 2>&1
set -efu
sv start pipewire # Hard depend on pipewire
. "$SVDIR"/../conf.sh # dbus session address

exec wireplumber
